# qlcplusmcp_raspi_service_pack/qlcplusmcp.service
[Unit]
Description=QLCPlus MCP HTTP Server
Documentation=https://github.com/infrafast/QLCPlus-MCP
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/QLCPlus-MCP
EnvironmentFile=/etc/qlcplusmcp.env
ExecStart=/usr/bin/npm run start:http
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
